# airootfs/etc/xdg/autostart/autorun.desktop
[Desktop Entry]
Encoding=UTF-8
Name=x11vnc
Comment=Init desktop interface
Exec=/home/tedy/autostart.sh
Terminal=true
Type=Application